# docker/root/etc/services.d/ferdi-server/run
#!/usr/bin/with-contenv bash

cd /app

setcap 'cap_net_bind_service=+ep' `which node`

# start server
echo " "
echo "**** Starting Ferdi-server ****"
exec s6-setuidgid abc adonis serve --dev